# etc/cron.daily/backupoffsite
#!/bin/sh

/root/backupoffsite.sh --quiet
