# etc/init.d/hostname_vps
#! /bin/bash
### BEGIN INIT INFO
# Provides: hostname_vps
# Required-Start:
# Required-Stop:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: set hostname.
# Description: Host Europe virtual server hostname setting script. Modified. -aj
### END INIT INFO

# /etc/init.d/hostname_vps
# https://kis.hosteurope.de/support/faq/index.php?cpid=15028

PATH=/sbin:/usr/sbin:/bin:/usr/bin

. /lib/init/vars.sh
. /lib/lsb/init-functions


#[ -f /etc/default/hostname_vps ] && HOSTNAME_VPS="$(cat /etc/default/hostname_vps)"
HOSTNAME_VPS=clyde.skgb.de


do_start () {
	[ -f /etc/hostname ] && echo "$HOSTNAME_VPS" > /etc/hostname
	[ -f /etc/HOSTNAME ] && echo "$HOSTNAME_VPS" > /etc/HOSTNAME
	[ -f /etc/mailname ] && echo "$HOSTNAME_VPS" > /etc/mailname
	[ -f /etc/sysconfig/network ] && sed -i.hebackup s/HOSTNAME=".*"/HOSTNAME=\"$HOSTNAME_VPS\"/ /etc/sysconfig/network
	/bin/hostname "$HOSTNAME_VPS"
	[ -f /etc/my.resolv.conf ] && cp /etc/my.resolv.conf /etc/resolv.conf
	[ -f /etc/my.hosts ] && cp /etc/my.hosts /etc/hosts
	
	# overwrite login message (the head of which is dynamic, generated by /etc/init.d/motd)
	#MOTD="`uname -o` `uname -r` (Debian `sed -e 's:jessie/sid:sid:' /etc/debian_version`) `uname -v`"
	MOTD="`uname -o` (Debian `sed -e 's:jessie/sid:sid:' /etc/debian_version`) `uname -v`"
	#echo "$MOTD" > /var/run/motd.dynamic  # /var/run/motd.dynamic is a symlink to /run/motd.dynamic
	> /run/motd.dynamic
	#echo >> /run/motd.dynamic
	echo "$MOTD" >> /run/motd.dynamic
}

do_status () {
	HOSTNAME=$(hostname)
	if [ "$HOSTNAME" == "$HOSTNAME_VPS" ]
	then
		return 0
	else
		return 4
	fi
}

case "$1" in
  start|"")
	do_start
	;;
  restart|reload|force-reload)
	echo "Error: argument '$1' not supported" >&2
	exit 3
	;;
  stop)
	# No-op
	;;
  status)
	do_status
	exit $?
	;;
  *)
	echo "Usage: hostname_vps start" >&2
	exit 3
	;;
esac

:
